# scripts/pi/service/meta.to.index.service
[Unit]
Description=meta.to.index
After=syslog.target rabbitmq-server.service

[Service]
EnvironmentFile=/etc/time.env
User=pi
WorkingDirectory=/var/lib/time
ExecStart=/usr/bin/java -jar /var/lib/time/meta.to.index.jar
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
